# hal/buildroot/linux/overlay/usr/lib/systemd/system/mount-host-certs.service
[Unit]
Description=Mount host certificates via 9p and configure mTLS
After=local-fs.target
Before=cube-agent.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mount-host-certs.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
